# runit/celery/run
#!/usr/bin/env bash

source /sandbox/venv/bin/activate

pushd /sandbox/sandbox

celery -A sandbox worker -l info -B -E